# cron.d/fetch.repos
PATH=$PATH:/usr/local/bin
*/1 * * * * root /etc/cron.d/task.sh task.fetch.repos.php > /tmp/log
