# cron.d/fetch.repos
PATH=$PATH:/usr/local/bin
*/1 * * * * root /etc/cron.d/task.sh task.fetch.repos.php
*/30 * * * * root /etc/cron.d/task.sh task.push.repos.php
