# cron.d/task.sh
#!/bin/sh

DIR=/home/lwch/public_github/bin # path to scripts path
TASK=$1

php $DIR/$TASK
